# bin/mesibo.service
[Unit]
Description=mesibo real-time communication server
Documentation=http://mesibo.com/documentation/
After=network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/var/run/mesibo.pid
ExecStart=/mesibo/bin/mesibo start
ExecReload=/mesibo/bin/mesibo reload
ExecStop=/mesibo/bin/mesibo stop
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
